# autostart/kleopatra.desktop
[Desktop Entry]
Type=Application
Name=Kleopatra
Exec=kleopatra
X-GNOME-Autostart-enabled=true
OnlyShowIn=UWSM;
